# systemd/uav-proc-airport.service
[Unit]
Description=UAV Robot Airport Process
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@PROJECT_ROOT@
EnvironmentFile=-/etc/default/uav_robot
ExecStart=/bin/bash -lc 'exec "@PROJECT_ROOT@/proc_airport/build/bin/proc_airport" ${UAV_PROC_AIRPORT_ARGS:-}'
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
