# systemd/uav-robotd.service
[Unit]
Description=UAV Robot Main Controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@PROJECT_ROOT@
EnvironmentFile=-/etc/default/uav_robot
ExecStart=/bin/bash -lc 'exec "@PROJECT_ROOT@/build/bin/uav_robotd" ${UAV_ROBOTD_ARGS:---listen-port 9001}'
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
